# manager/supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/usr/local/bro/bin/init_cluster_broctl.sh

[program:sshd]
command=/usr/sbin/sshd -D