# manager/supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/bro/etc/cluster_config/start-command.sh

[program:sshd]
command=/usr/sbin/sshd -D